# system/user-units/rangeos-session.target
[Unit]
Description=RangeOS AI Session Control Target
AllowIsolate=yes
Wants=rangeos-shell.service

[Install]
WantedBy=default.target
